# video-worker/crontab
*/5 * * * * node /root/index.js >> /var/log/cron.log 2>&1
